# varnish.conf
[supervisord]
nodaemon=true

[program:varnishd]
command=/usr/sbin/varnishd -f /etc/varnish/default.vcl  -a :%(ENV_VCL_PORT)s -s malloc,%(ENV_MAX_MEMORY)s -F
autorestart=true
