# varnish.conf
[supervisord]
nodaemon=true

[program:varnishd]
command=/usr/sbin/varnishd -f /etc/varnish/default.vcl  -a ${HTTP_PORT} -F
autorestart=true
